# destriping_scripts_rotational_buffer/Destripe_C_Rotational_buffer_L1B_SDR.cron
0 */6 * * * /data/home001/mikelson/destriping_scripts_rotational_buffer/script_rotational_buffer_destripe >> /data/home001/mikelson/destriping_scripts_rotational_buffer/script_rotational_buffer_destripe.log
